# deployment/rc.local
#!/bin/sh -e

# This script runs each time the machine reboots
# place it in /etc and be sure that it is executable

# Map External Ports
iptables -t nat -A PREROUTING -i eth0 -p tcp --dport 80 -j REDIRECT --to-port 3010

# Start Mongo Database Daemon
mongod > /dev/null 2>&1 &

# Wait until mongo has started
sleep 10s

# Start webserver
runuser -l ubuntu -c 'forever start /home/ubuntu/image-carousel-component/server/index.js'